# contrib/cron/ocsinventory-agent.cron
#!/bin/bash
NAME=ocsinventory-agent

exec >>/var/log/$NAME/$NAME.log 2>&1 

[ -f   /etc/sysconfig/$NAME ] || exit 0
source /etc/sysconfig/$NAME
export PATH

i=0
while [ $i -lt ${#OCSMODE[*]} ]
do
	if [ ${OCSMODE[$i]:-none} == cron ]; then
		OPTS=
		if [ ! -z "${OCSPAUSE[$i]}" ]; then
			OPTS="--wait ${OCSPAUSE[$i]}"
		fi

		if [ ! -z "${OCSTAG[$i]}" ]; then
		        OPTS="$OPTS --tag=${OCSTAG[$i]}"
		fi

		if [ "z${OCSSERVER[$i]}" = 'zlocal' ]; then
	        	# Local inventory
	        	OPTS="$OPTS --local=/var/lib/$NAME"

		elif [ ! -z "${OCSSERVER[$i]}" ]; then
	        	# Remote inventory
		        OPTS="$OPTS --lazy --nolocal --server=${OCSSERVER[$i]}"
		        if [ ! -z "${OCSPROXYSERVER[$i]}" ]; then
		        	OPTS="$OPTS --proxy=${OCSPROXYSERVER[$i]}"
		        fi
		fi

		echo "[$(date '+%c')] Running $NAME $OPTS"
		/usr/sbin/$NAME  $OPTS
	fi
	((i++))
done
echo "[$(date '+%c')] End of cron job ($PATH)"

